# systemd/hermes-full-backup-daily.timer
[Unit]
Description=Create a daily full Hermes backup

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15m
Unit=hermes-full-backup@daily.service

[Install]
WantedBy=timers.target
